# ansible/bastion/roles/clair/templates/container-clair-clair.service.j2
# container-clair-clair.service
# autogenerated by Podman 1.9.3
# Thu Sep 10 11:58:41 CEST 2020

[Unit]
Description=Podman container-clair-clair.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=container-clair-postgres.service
Requires=container-clair-postgres.service


[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
ExecStartPre=/usr/bin/rm -f %t/%n-pid %t/%n-cid
ExecStart=/usr/bin/podman run \
    --conmon-pidfile %t/%n-pid \
    --cidfile %t/%n-cid \
    --cgroups=no-conmon \
    -d \
    --rm \
    --ip {{ okd_lab_clair_ip }} \
    -p 6060:6060 \
    -p 6061:6061 \
    --name=clair-clair \
    --privileged=true \
    -v /var/opt/clair/config:/clair/config:Z \
    -v /var/opt/clair/config:/etc/pki/ca-trust/source/anchors:Z \
    {{ okd_lab_quay_clair_image }}
ExecStop=/usr/bin/podman stop --ignore --cidfile %t/%n-cid -t 10
ExecStopPost=/usr/bin/podman rm --ignore -f --cidfile %t/%n-cid
PIDFile=%t/%n-pid
KillMode=none
Type=forking

[Install]
WantedBy=multi-user.target default.target
