# ansible/bastion/roles/artifactory/templates/container-artifactory-postgres.service.j2
# container-artifactory-postgres.service
# autogenerated by Podman 1.9.3
# Thu Sep 10 08:14:09 CEST 2020

[Unit]
Description=Podman container-artifactory-postgres.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n

Restart=always
ExecStartPre=/usr/bin/rm -f %t/%n-pid %t/%n-cid
ExecStart=/usr/bin/podman run \
    --conmon-pidfile %t/%n-pid \
    --cidfile %t/%n-cid \
    --cgroups=no-conmon \
    -d \
    --name artifactory-postgres \
    --ip {{ okd_lab_artifactory_postgres_ip }} \
    -e POSTGRES_USER=artifactory \
    -e POSTGRES_PASSWORD=password \
    -e POSTGRES_DB=artifactory \
    -v /var/opt/artifactory/postgres:/var/lib/postgresql/data:Z \
    {{ okd_lab_artifactory_postgres_image }}
ExecStop=/usr/bin/podman stop --ignore --cidfile %t/%n-cid -t 10
ExecStopPost=/usr/bin/podman rm --ignore -f --cidfile %t/%n-cid
PIDFile=%t/%n-pid
KillMode=none
Type=forking

[Install]
WantedBy=multi-user.target default.target

