# ansible/bastion/roles/artifactory/templates/container-artifactory-artifactory.service.j2
# container-artifactory-artifactory.service
# autogenerated by Podman 1.9.3
# Thu Sep 10 11:58:41 CEST 2020

[Unit]
Description=Podman container-artifactory-artifactory.service
Documentation=man:podman-generate-systemd(1)
After=container-artifactory-postgres.service
Requires=container-artifactory-postgres.service


[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n

Restart=on-failure
ExecStartPre=/usr/bin/rm -f %t/%n-pid %t/%n-cid
ExecStart=/usr/bin/podman run \
    --conmon-pidfile %t/%n-pid \
    --cidfile %t/%n-cid \
    --cgroups=no-conmon \
    --rm \
    -d \
    --ip {{ okd_lab_artifactory_ip }} \
    --name=artifactory-artifactory \
    -v /var/opt/artifactory/var:/var/opt/jfrog/artifactory:Z \
    -v /etc/localtime:/etc/localtime:ro \
    -e EXTRA_JAVA_OPTIONS="{{ okd_lab_artifactory_java_opts }}" \
    -e JF_ROUTER_ENTRYPOINTS_EXTERNALPORT=8082 \
    {{ okd_lab_artifactory_image }}
ExecStop=/usr/bin/podman stop --ignore --cidfile %t/%n-cid -t 10
ExecStopPost=/usr/bin/podman rm --ignore -f --cidfile %t/%n-cid
PIDFile=%t/%n-pid
KillMode=none
Type=forking

[Install]
WantedBy=multi-user.target default.target
