# ansible/bastion/roles/gitlab/templates/container-gitlab.service.j2
# container-gitlab.service
# autogenerated by Podman 1.9.3
# Thu Sep 10 09:23:05 CEST 2020

[Unit]
Description=Podman container-gitlab.service
Documentation=man:podman-generate-systemd(1)
After=network-online.target
Wants=network.target

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
TimeoutSec=900

Restart=on-failure
ExecStartPre=/usr/bin/rm -f %t/%n-pid %t/%n-cid

# Check health-checks, disabled

ExecStart=/usr/bin/podman run \
    --conmon-pidfile %t/%n-pid \
    --cidfile %t/%n-cid \
    --cgroups=no-conmon \
    --health-start-period 2m3s \
    --health-timeout 6m \
    --health-interval disable \
    --health-cmd none \
    --rm \
    -d \
    --name gitlab \
    --ip {{ okd_lab_gitlab_ip }} \
    --publish 2222:22 \
    --volume /var/opt/gitlab/config:/etc/gitlab:Z \
    --volume /var/opt/gitlab/logs:/var/log/gitlab:Z \
    --volume /var/opt/gitlab/data:/var/opt/gitlab:Z \
    --volume /var/opt/gitlab/registry:/var/opt/gitlab/gitlab-rails/shared/registry:Z \
    -e GITLAB_ROOT_PASSWORD={{ okd_lab_gitlab_root_password }} \
    -e CHEF_FIPS= \
    -e GITLAB_ROOT_EMAIL={{ okd_lab_gitlab_root_email }} \
    {{ okd_lab_gitlab_container_image }}
ExecStop=/usr/bin/podman stop --ignore --cidfile %t/%n-cid -t 10
ExecStopPost=/usr/bin/podman rm --ignore -f --cidfile %t/%n-cid
PIDFile=%t/%n-pid
KillMode=none
Type=forking

[Install]
WantedBy=multi-user.target default.target

