# ansible/bastion/roles/quay/templates/container-quay-redis.service.j2
# container-quay-redis.service
# autogenerated by Podman 1.9.3
# Thu Sep 10 08:59:49 CEST 2020

[Unit]
Description=Podman container-quay-redis.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
ExecStartPre=/usr/bin/rm -f %t/%n-pid %t/%n-cid
ExecStart=/usr/bin/podman run \
    --conmon-pidfile %t/%n-pid \
    --cidfile %t/%n-cid \
    --cgroups=no-conmon \
    -d \
    --name quay-redis \
    --ip {{ okd_lab_quay_redis_ip }} \
    {{ okd_lab_quay_redis_image }} --requirepass strongpassword
ExecStop=/usr/bin/podman stop --ignore --cidfile %t/%n-cid -t 10
ExecStopPost=/usr/bin/podman rm --ignore -f --cidfile %t/%n-cid
PIDFile=%t/%n-pid
KillMode=none
Type=forking

[Install]
WantedBy=multi-user.target default.target

