# ansible/bastion/roles/quay/templates/container-quay-postgres.service.j2
# container-quay-postgres.service
# autogenerated by Podman 1.9.3
# Thu Sep 10 08:14:09 CEST 2020

[Unit]
Description=Podman container-quay-postgres.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
TimeoutSec=900

Restart=always
ExecStartPre=/usr/bin/rm -f %t/%n-pid %t/%n-cid
ExecStart=/usr/bin/podman run \
    --conmon-pidfile %t/%n-pid \
    --cidfile %t/%n-cid \
    --cgroups=no-conmon \
    --rm \
    -d \
    --name quay-postgres \
    --ip {{ okd_lab_quay_postgres_ip }} \
    -e POSTGRES_USER={{ okd_lab_quay_postgres_user }} \
    -e POSTGRES_PASSWORD={{ okd_lab_quay_postgres_password }} \
    -e POSTGRES_DB=quay \
    -v /var/opt/quay/postgres:/var/lib/postgresql/data:Z \
    {{ okd_lab_quay_postgres_image }}
ExecStop=/usr/bin/podman stop --ignore --cidfile %t/%n-cid -t 10
ExecStopPost=/usr/bin/podman rm --ignore -f --cidfile %t/%n-cid
PIDFile=%t/%n-pid
KillMode=none
Type=forking

[Install]
WantedBy=multi-user.target default.target

