# ansible/bastion/roles/quay/templates/container-quay-quay.service.j2
# container-quay-quay.service
# autogenerated by Podman 1.9.3
# Thu Sep 10 11:58:41 CEST 2020

[Unit]
Description=Podman container-quay-quay.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=container-quay-postgres.service
After=container-quay-redis.service
Requires=container-quay-postgres.service
Requires=container-quay-redis.service


[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
ExecStartPre=/usr/bin/rm -f %t/%n-pid %t/%n-cid
ExecStart=/usr/bin/podman run \
    --conmon-pidfile %t/%n-pid \
    --cidfile %t/%n-cid \
    --cgroups=no-conmon \
    --rm \
    -d \
    --ip {{ okd_lab_quay_ip }} \
    --name=quay-quay \
    --privileged=true \
    -v /var/opt/quay/config:/conf/stack:Z \
    -v /var/opt/quay/storage:/datastorage:Z \
    {{ okd_lab_quay_image }}
ExecStop=/usr/bin/podman stop --ignore --cidfile %t/%n-cid -t 10
ExecStopPost=/usr/bin/podman rm --ignore -f --cidfile %t/%n-cid
PIDFile=%t/%n-pid
KillMode=none
Type=forking

[Install]
WantedBy=multi-user.target default.target
