# ansible/bastion/roles/quay/templates/container-quay-config.service.j2
# container-quay-config.service
# autogenerated by Podman 1.9.3
# Thu Sep 10 09:23:05 CEST 2020

[Unit]
Description=Podman container-quay-config.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
Requires=container-quay-postgres.service
After=container-quay-postgres.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
ExecStartPre=/usr/bin/rm -f %t/%n-pid %t/%n-cid
ExecStart=/usr/bin/podman run \
    --conmon-pidfile %t/%n-pid \
    --cidfile %t/%n-cid \
    --cgroups=no-conmon \
    -d \
    --privileged=true \
    --rm \
    -it \
    --ip {{ okd_lab_quay_config_ip }} \
    --name quay-config \
    -p 8081:8443 \
    {{ okd_lab_quay_image }} config {{ okd_lab_quay_config_password }}
ExecStop=/usr/bin/podman stop --ignore --cidfile %t/%n-cid -t 10
ExecStopPost=/usr/bin/podman rm --ignore -f --cidfile %t/%n-cid
PIDFile=%t/%n-pid
KillMode=none
Type=forking

[Install]
WantedBy=multi-user.target default.target

